# scripts/com.crypto-info.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.crypto-info</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/crypto-info/scripts/wrapper.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/opt/crypto-info</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>StandardOutPath</key>
    <string>/opt/crypto-info/backend/logs/launchd-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/opt/crypto-info/backend/logs/launchd-stderr.log</string>
</dict>
</plist>
